# deploy/lounge-bot.service
[Unit]
Description=Discord Lounge Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/.openclaw/workspace/discord-bot
ExecStart=/home/ec2-user/.openclaw/workspace/discord-bot/venv/bin/python3 lounge_bot.py
EnvironmentFile=/home/ec2-user/.openclaw/workspace/discord-bot/.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
